# systemd/user/wayland-session-bindpid@.service.in
[Unit]
Description=Bind graphical session to PID %i
Documentation=man:uwsm(1)
OnSuccess=wayland-session-shutdown.target
OnSuccessJobMode=replace-irreversibly
OnFailure=wayland-session-shutdown.target
OnFailureJobMode=replace-irreversibly
Conflicts=wayland-session-shutdown.target
Before=wayland-session-shutdown.target
CollectMode=inactive-or-failed
[Service]
Type=exec
ExecStart=/bin/sh -c "if command -v waitpid >/dev/null; then exec waitpid -e %i; else exec @BIN_PATH@ aux waitpid %i; fi" waitpid
Restart=no
SyslogIdentifier=uwsm_bindpid
Slice=background.slice
