# systemd/user/wayland-wm@.service.in
[Unit]
Description=Main service for %I
Documentation=man:uwsm(1)
Requires=wayland-session-pre@%i.target
BindsTo=wayland-session@%i.target
Before=wayland-session@%i.target graphical-session.target
PropagatesStopTo=wayland-session@%i.target graphical-session.target
After=wayland-session-pre@%i.target graphical-session-pre.target
Wants=wayland-session-envelope@%i.target
OnSuccess=wayland-session-shutdown.target
OnSuccessJobMode=replace-irreversibly
OnFailure=wayland-session-shutdown.target
OnFailureJobMode=replace-irreversibly
Conflicts=wayland-session-shutdown.target
Before=wayland-session-shutdown.target
CollectMode=inactive-or-failed
[Service]
# awaits for ready state notification from compositor or child
# should be issued by '@BIN_NAME@ finalize' for example
Type=notify
NotifyAccess=all
ExecStart=@BIN_PATH@ aux exec -- %I
Restart=no
EnvironmentFile=-%t/@BIN_NAME@/env_session.conf
TimeoutStartSec=30
TimeoutStopSec=10
SyslogIdentifier=@BIN_NAME@_%I
Slice=session.slice
